# scripts/launchd/com.jeongsik.wkbl-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.jeongsik.wkbl-refresh</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-lc</string>
      <string>/Users/dancer/me/workspace/yousleepwhen/wkbl/scripts/run_wkbl_refresh.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/dancer/me/workspace/yousleepwhen/wkbl</string>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>4</integer>
      <key>Minute</key>
      <integer>10</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/dancer/me/logs/wkbl_refresh.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dancer/me/logs/wkbl_refresh.launchd.log</string>
  </dict>
</plist>
