# ocaml/scripts/com.wkbl.schedule-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wkbl.schedule-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/dancer/me/workspace/yousleepwhen/wkbl/ocaml/scripts/sync-schedule.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <!-- Run daily at 00:30 (after games end) -->
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/dancer/me/workspace/yousleepwhen/wkbl/ocaml/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/dancer/me/workspace/yousleepwhen/wkbl/ocaml/logs/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/dancer/.opam/default/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/dancer/me/workspace/yousleepwhen/wkbl/ocaml</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
